# assignement04/75-keyboard.rules
SUBSYSTEM=="usb", ATTR{product}=="Storage Media", RUN+="/sbin/insmod /home/little-penguin/rendu/assignement04/main.ko"
