# assignement04/75-keyboard.rules
SUBSYSTEM=="usb", ATTR{product}=="Apple Keyboard", RUN+="/sbin/insmod /home/ariard/little-penguin1/assignement04/main.ko"
